# deploy/audit-guru-worker.service
[Unit]
Description=Audit Guru - SQS queue worker (runs the OCR -> Validation -> Audit pipeline)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/audit-guru
Environment=AWS_REGION=us-east-1
ExecStart=/home/ec2-user/audit-guru/venv/bin/python queue_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
